# deploy/campus-keepalive.service.example
[Unit]
Description=Campus Dr.COM auto connect and keepalive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/campus-keepalive.env
ExecStart=/usr/local/bin/campus-drcom-keepalive --env-file /etc/campus-keepalive.env --gateway-cache-file /var/lib/campus-keepalive/default.gateway
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
